# rootfs/icarus-grimoire.timer
[Unit]
Description=Run Grimoire Keeper every hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target